# data/appgrid.conf
# appgrid - monitors the states of packages to keep the appgrid db up to date

description	"monitors installation states"

start on runlevel [12345]
stop on runlevel [!12345]

exec /usr/bin/appgrid --monitor-state
